# deploy/rustspray.service
[Unit]
Description=Rust-Spray Agricultural Spraying Pipeline
Documentation=https://github.com/cropcrusaders/Rust-Spray
After=network.target

[Service]
Type=simple

# Camera capture piped into the spray pipeline binary.
# The rustspray-camera script reads /etc/rustspray/config.toml
# to determine camera backend, resolution, and frame rate.
ExecStart=/bin/sh -c '/usr/local/bin/rustspray-camera | /usr/local/bin/rustspray --config /etc/rustspray/config.toml'

# Always restart: a camera drop-out makes rustspray exit cleanly on EOF,
# and the sprayer must come back on its own in the field.
Restart=always
RestartSec=5

# GPIO access requires root.
User=root

# Logging goes to the journal (view with: journalctl -u rustspray -f)
StandardOutput=journal
StandardError=journal

Environment=RUST_LOG=info

# Hardening. GPIO access goes through /dev/gpiomem* (rppal), which is
# unaffected by ProtectSystem. Do NOT add ReadWritePaths=/sys/class/gpio —
# that path no longer exists on Raspberry Pi OS Bookworm and a missing
# ReadWritePaths entry makes the unit fail to start.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
